# systemd/zdd-detector.service
[Unit]
Description=Zero-Day Detector System
Documentation=https://github.com/zdd-security/detector
After=network.target
Wants=network.target

[Service]
Type=simple
User=zdd-detector
Group=zdd-detector
ExecStart=/usr/local/bin/python3 /usr/local/bin/zdd-detector --config /etc/zdd/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/zdd /var/lib/zdd /tmp

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Capabilities needed for packet capture
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

# Environment
Environment=PYTHONPATH=/usr/local/lib/zdd
Environment=ZDD_CONFIG_PATH=/etc/zdd

[Install]
WantedBy=multi-user.target